# debian/etc/init/astrobox.conf
# astrobox - astrobox sofware
#
# Upstart configuration file for the AstroBox Software

description "AstroBox(tm) software"

start on (filesystem and started networking)

respawn

pre-start script

        mount -t vfat /dev/mmcblk0p1 /sdcard -o rw,uid=ubuntu,gid=ubuntu

        if ! /usr/bin/network-monitor check; then
                /usr/bin/network-monitor hotspot
        fi

end script

post-stop script

        umount /sdcard

end script

exec python -OO /AstroBox/run --config /etc/astrobox/config.yaml --host 127.0.0.1