# install/init/rc.local
### AstroBox rc.local ### (Automatically generated)

if ! /usr/bin/network-monitor check; then
	/usr/bin/network-monitor hotspot
fi

###
